# upstart_configurations/dnsmasq-consul-template.conf
description "Start consul-template for reloading dnsmasq config on boot"
author      "Suman Mukherjee"

start on containers-ready
stop on shutdown

script
  exec consul-template -consul 172.20.20.11:8500 -wait 2s -retry 5s -pid-file "/var/run/consul-template.pid" -template "/vagrant/docker/dnsmasq/2.68/configuration/hosts.ctmpl:/etc/dnsmasq.hosts:docker stop dnsmasqServer && docker start dnsmasqServer"
end script

pre-start script
  rm /etc/dnsmasq.hosts
  touch /etc/dnsmasq.hosts
  echo "Starting the dnsmasq consul template daemon now" >> /var/log/consul-template.log
end script

pre-stop script
  rm /var/run/consul-template.pid
  echo "Stopping the dnsmasq consul template daemon now" >> /var/log/consul-template.log
end script
